# deploy/auto-deploy.timer
[Unit]
Description=LIMS Auto-Deploy Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
